# scripts/systemd/nira-model.service
[Unit]
Description=NiraOS Model Manager
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=nira-ai
Group=nira-ipc
ExecStart=/usr/bin/nira-model
Restart=on-failure
RestartSec=5s
Environment="RUST_LOG=info"
NoNewPrivileges=yes
PrivateTmp=yes
ProtectHome=read-only
ProtectSystem=strict
ReadWritePaths=/var/lib/niraos/models /var/lib/niraos/ai /run/niraos
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6

[Install]
WantedBy=multi-user.target
