# mkosi/mkosi.extra/etc/systemd/system/NetworkManager.service.d/10-niraos-install.conf
[Install]
# Keep NetworkManager's dispatcher integration, but do not let the vendor
# unit's Also= pull wait-online back in after the NiraOS preset disables it.
Also=
Also=NetworkManager-dispatcher.service
