# systemd/user/nira-shell.service
[Unit]
Description=Nira Shell
Requires=nira-compositor.service
After=nira-compositor.service
# Not enabled via systemd — launched manually by start-nira-session
# after the Wayland compositor socket is ready.

[Service]
ExecStart=/usr/bin/nira-shell
Restart=on-failure
RestartSec=5s
